# backend/scripts/supervisord-new.conf
[supervisord]
nodaemon=true
user=app
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/tmp
silent=true

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700
chown=app:app

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# PHP-FPM Process
[program:php-fpm]
command=php-fpm --nodaemonize
user=app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=1

# Nginx Process
[program:nginx]
command=nginx -g 'daemon off; error_log stderr info;'
user=app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=2